# user/pulseaudio.service.d/jack.conf
[Service]
ExecStartPre=/usr/bin/jack_wait -w
