# user/pulseaudio.service.d/realtime.conf
[Service]
LimitRTPRIO=20
LimitRTTIME=infinity
#ExecStartPost=/bin/sh -c 'chrt -v -r -p 20 `pidof pulseaudio`'