# user/pulseaudio.service.d/realtime.conf
[Service]
LimitRTPRIO=20
LimitRTTIME=infinity
ExecStartPre=/usr/bin/jack_wait -w
