# user/pulseaudio.service.d/realtime.conf
[Service]
LimitRTPRIO=20
LimitRTTIME=infinity
