# user/jack-alsa-out@jabra.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-sound-Jabra_BIZ_2400_USB.device
After=sys-subsystem-sound-Jabra_BIZ_2400_USB.device

[Install]
WantedBy=sys-subsystem-sound-Jabra_BIZ_2400_USB.device
