# user/jack-alsa-in@X18XR18.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-usb-behringer_x18xr18.device
After=sys-subsystem-usb-behringer_x18xr18.device

[Install]
WantedBy=sys-subsystem-usb-behringer_x18xr18.device
