# user/jack-alsa-in@X18XR18.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-sound-X18_XR18.device
After=sys-subsystem-sound-X18_XR18.device

[Install]
WantedBy=sys-subsystem-sound-X18_XR18.device
