# user/jack-alsa-in@jabra.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-usb-jabra.device
After=sys-subsystem-usb-jabra.device

[Install]
WantedBy=sys-subsystem-usb-jabra.device
