# user/alsa-setup.service
[Unit]
Description=ALSA card(s) setup
Wants=environment.target
After=environment.target

[Service]
Type=oneshot
ExecStart=%E/systemd/bin/alsa-setup

[Install]
WantedBy=audio-setup.target
