# user/environment.target
[Unit]
Description=Environment Setup for the User Session
Documentation=man:systemd.special(7)
Requires=basic.target
After=basic.target

