# user/environment.target
[Unit]
Description=Setup environment for user services.
Wants=basic.target
After=basic.target
