# user/environment.target
[Unit]
Description=Environment for user services
Wants=basic.target
After=basic.target
