# user/X11@.target
[Unit]
Description=X11 Display %I
BindsTo=graphical-session.target

[Install]
WantedBy=xsession@%i.target
