# user/xsession@.target
[Unit]
Description=X11 Session on Display %I
Wants=X11@.target window-manager@.service xcompositor@.service
