# user/jack-alsa-in@gt1000.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-sound-GT\x2d1000.device
After=sys-subsystem-sound-GT\x2d1000.device

[Install]
WantedBy=sys-subsystem-sound-GT\x2d1000.device
