# user/jack-alsa-in@gt1000.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-usb-boss_gt1000.device
After=sys-subsystem-usb-boss_gt1000.device

[Install]
WantedBy=sys-subsystem-usb-boss_gt1000.device
