# user/mpd.service.d/jack.conf
[Unit]
Requires=jack.service
After=jack.service
