# user/mpd.service.d/network.conf
[Unit]
Requires=network.target
After=network.target
