# user/mpd.service.d/pulseaudio.conf
[Unit]
Requires=pulseaudio.service
After=pulseaudio.service
Requires=network.target
After=network.target
