# user/mpd.service.d/pulseaudio.conf
[Unit]
Requires=pulseaudio.service
After=pulseaudio.service
