# user/mpd.service.d/pulseaudio.conf
[Unit]
BindsTo=pulseaudio.service
After=pulseaudio.service
